# src/test/resources/hudson/plugins/active_directory/docker/TheFlintstonesTest/TheFlintstones/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:bind9]
command=/usr/sbin/named -c /etc/bind/named.conf -u bind -f

[program:kerberos]
command=/usr/sbin/krb5kdc -n

[program:samba]
command=/usr/sbin/samba -i

[program:sssd]
command=/usr/sbin/sssd -i

[program:syslog]
command=/usr/sbin/rsyslogd -n

[program:custom]
command=/usr/local/bin/custom.sh